# confs/supervisor/conf.d/kibana.conf
[program:kibana]
command=/home/pi/Downloads/kibana/bin/kibana serve
autostart=true
autorestart=unexpected
numprocs=1
stdout_logfile=/tmp/kibana.log
stderr_logfile=/tmp/kibana.err
